# config/supervisor/conf.d/mariadb.conf
[program:mariadb]
command=/usr/bin/mysqld --user=mysql --datadir=/var/lib/mysql --skip-networking=0 --bind-address=127.0.0.1
user=mysql
autostart=true
autorestart=true
priority=1
